# mirror/crontab
# Arch Linux mirrors
* * * * * /syncrepo.sh

# Status heartbeats
*/1 * * * * wget --spider "$UP_HEARTBEAT_URL" >/dev/null 2>&1